# demo/infra/openclaw-demo.service
[Unit]
Description=OpenClaw Demo Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/openclaw gateway run
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/home/openclaw/.openclaw/.env

[Install]
WantedBy=default.target
